# launcher/scala-akka-http-server-test.service
[Unit]
Description= start and stop scala-akka-http-server for Ubuntu (systemd)
After=network.target

[Service]
Type=simple
User=vvasuki
WorkingDirectory=/home/vvasuki/scala-akka-http-server/
ExecStart=/usr/bin/sbt run

[Install]
WantedBy=multi-user.target